# roles/naive_proxy/templates/acme-renew.timer.j2
[Unit]
Description=Daily ACME certificate renewal for naive-proxy

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
