# roles/iptables/templates/mini-pig-firewall.service.j2
{{ ansible_managed | comment }}
[Unit]
# Type=oneshot, no RemainAfterExit: each start re-asserts the ruleset
# (idempotent; "inactive (dead)" between runs is expected).
# After=network-online, NOT Before=network-pre: re-applying mid NIC bring-up
# (`nft -f` atomically delete+recreates the mpig_* tables) triggers a kernel
# nftables lockup — suspected dynamic rate-limit set teardown vs packet-path
# race — that hard-hangs boot on cloud-init hosts that block until
# network-online (e.g. Scaleway Elastic Metal).
Description=Mini Pig firewall (idempotent apply + drift check)
Wants=network-online.target
After=network-online.target
ConditionPathExists=/usr/local/sbin/mini-pig-firewall-apply

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mini-pig-firewall-apply

[Install]
WantedBy=multi-user.target
