# roles/iptables/templates/mpig-randomized-snat.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=Mini Pig randomized SNAT (native nft chain at priority srcnat - 10)
# After=network-online, NOT Before=network-pre: same kernel nftables bring-up
# lockup as mini-pig-firewall.service (re-apply mid NIC bring-up freezes boot
# on cloud-init hosts). After=nftables.service so the SNAT chain lands after
# the base ruleset exists.
Wants=network-online.target
After=network-online.target nftables.service
ConditionPathExists=/etc/nftables.d/mpig-randomized-snat.conf

[Service]
# Type=oneshot without RemainAfterExit — mirrors mini-pig-firewall.service.
# Each start atomically re-applies the chain. The atomic delete-then-create
# lives INSIDE the conf file (`add table; delete table; table {…}` idiom,
# see AGENTS.md pitfall P4), so a bare `nft -f` is one kernel transaction
# — no shell wrapper, no race window. Periodic re-application is driven
# by mpig-randomized-snat.timer. No PartOf= or ReloadPropagatedFrom=
# coupling to nftables.service: the two units own independent nft tables
# and stay decoupled.
Type=oneshot
ExecStart=/usr/sbin/nft -f /etc/nftables.d/mpig-randomized-snat.conf

[Install]
WantedBy=multi-user.target
